# deploy_tools/lawn-care-planner.conf
[program:lawn-care-planner]
command =sh /home/lcp/bin/gunicorn_start
user=lcp
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/lcp/logs/gunicorn-error.log
